# com.liorbar.greenlight.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.liorbar.greenlight</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/liorbar/.local/share/uv/python/cpython-3.11.14-macos-aarch64-none/bin/python3.11</string>
        <string>/Users/liorbar/Documents/all_projects/greenlight/greenlight_app.py</string>
    </array>

    <!-- Start at login. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart only on abnormal exit (crash). A clean quit (double-click)
         exits 0 and is respected; the single-instance lock means a duplicate
         launch exits 0 too and is likewise not restarted. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- GUI session only. -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>StandardOutPath</key>
    <string>/Users/liorbar/Documents/all_projects/greenlight/app.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/liorbar/Documents/all_projects/greenlight/app.log</string>
</dict>
</plist>
